# scripts/tellar.service
[Unit]
Description=Tellar Cyber Steward
After=network.target

[Service]
Type=simple
# The setup script will replace {{GUILD_PATH}} with your actual guild directory
WorkingDirectory={{GUILD_PATH}}
ExecStart=%h/.cargo/bin/tellar --guild {{GUILD_PATH}}
Restart=always
RestartSec=10
# Standard output/error logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
